# webapp.conf
[program:webapp]
command=/usr/bin/gunicorn app:app
directory=/home/pi/garage-zero
autostart=true
autorestart=true
startretries=3
stderr_logfile=/home/pi/garage-zero/logs/webapp.err.log
stdout_logfile=/home/pi/garage-zero/logs/webapp.out.log
user=pi
